# config/graccsumperiodicyearly.timer
[Unit]
Description=Run graccsumperiodic every day and on boot

[Timer]
OnCalendar=1day
Persistent=true

[Install]
WantedBy=timers.target